# ansible/roles/gw/templates/renew.timer
[Unit]
Description=Run renew certificates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
