# ansible/roles/rasp/templates/renew.service
[Unit]
Description=Certificates
PartOf=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/bitwarden
ExecStart=/usr/bin/docker run --rm --name certbot \
	-v "/etc/letsencrypt:/etc/letsencrypt" \
	-v "/opt/bitwarden:/opt/bitwarden" \
	-v "/var/lib/letsencrypt:/var/lib/letsencrypt" \
	certbot/dns-ovh renew
ExecStartPost=docker kill -s SIGHUP lb
