# roles/certs/templates/renew_rfc2136.timer
[Unit]
Description=Run RFC2136 renew certificates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
